# systemd/onex-dispatch-claim-sweeper.service
# SPDX-FileCopyrightText: 2025 OmniNode.ai Inc.
# SPDX-License-Identifier: MIT
[Unit]
Description=ONEX Dispatch Claim Sweeper — reap expired claim files
After=network.target

[Service]
Type=oneshot
Environment=ONEX_STATE_DIR=%h/.onex
ExecStart=/bin/sh -c 'cd "%h" && uv run python -m onex_change_control.dispatch_claims.sweeper'
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
