# systemd/onex-dispatch-claim-sweeper.timer
# SPDX-FileCopyrightText: 2025 OmniNode.ai Inc.
# SPDX-License-Identifier: MIT
[Unit]
Description=ONEX Dispatch Claim Sweeper timer — every 5 minutes
Requires=onex-dispatch-claim-sweeper.service

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
Unit=onex-dispatch-claim-sweeper.service

[Install]
WantedBy=timers.target
